# Cluster_Automotivo/cluster-automotivo.service
[Unit]
Description=Cluster Automotivo - Painel de Instrumentos Web
Documentation=https://github.com/marcosgt96/Rep_Cluster
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=%h/Rep_Cluster/Cluster_Automotivo
Environment="PATH=%h/Rep_Cluster/Cluster_Automotivo/cluster_env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/usr/bin/env bash -lc 'cd %h/Rep_Cluster/Cluster_Automotivo && exec ./run_cluster.sh'
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cluster-automotivo

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=%h/Rep_Cluster/Cluster_Automotivo/logs

[Install]
WantedBy=multi-user.target
