# Cluster_Automotivo/cluster-automotivo.service
[Unit]
Description=Cluster Automotivo - Painel de Instrumentos Web
Documentation=https://github.com/marcosgt96/Rep_Cluster
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory={{APP_DIR}}
Environment="PATH={{APP_DIR}}/cluster_env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/bin/bash -lc 'cd {{APP_DIR}} && ./run_cluster.sh'
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cluster-automotivo

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{APP_DIR}}/logs

[Install]
WantedBy=multi-user.target
